# files/rsyslogd.rc.OpenBSD
#!/bin/sh

daemon="/usr/local/sbin/rsyslogd"

. /etc/rc.d/rc.subr

rc_cmd $1
